# config/supervisord/CampusQueriesCard.conf
[program:CampusQueriesCard]
command=/usr/local/bin/supervisor -e 'node,js,css,html' -w /var/www/query node /var/www/query/app.js
autorestart=false
user=root
